# deployment/minecraft.service
[Unit]
Description=Minecraft Server
After=network.target

[Service]
User=nakanishi_ryot
WorkingDirectory=/home/nakanishi_ryot/minecraft_server
# Adjust memory (Xms/Xmx) based on GCE instance type
ExecStart=/usr/lib/jvm/java-21-openjdk-amd64/bin/java -Xms2G -Xmx6G -jar server.jar --nogui
Restart=always
RestartSec=20

[Install]
WantedBy=multi-user.target